# skel/etc/profile.d/ade_blackdrones_env.sh
#!/usr/bin/env bash


source /opt/ros/eloquent/setup.sh

cd
